# debian/istatd-server.prerm
#!/bin/sh
set -e

if [ "$1" != "remove" -a "$1" != "purge" ]; then
        exit 0
fi

if [ -x "/etc/init.d/istatd-server" ]; then
        update-rc.d istatd-server defaults
        if [ -x /usr/sbin/invoke-rc.d ]; then
                invoke-rc.d istatd-server stop || true
        else
                /etc/init.d/istatd-server stop || true
        fi
fi

#DEBHELPER#

exit 0
